# nexus.conf
[supervisord]
nodaemon=true

[program:nexus]
user=nexus
command=/bin/bash -c "/opt/supervisor/nexus_supervisor"
redirect_stderr=true

